# conf.d/bookmark-processor.conf
[program:bookmark-processor]
command=/Users/elenafaillace/.local/bin/uv run uvicorn app.api:app --port 8484
directory=/Users/elenafaillace/Documents/all_code/bookmark-content-processor
autostart=true
autorestart=true
stdout_logfile=/Users/elenafaillace/Library/Logs/supervisor/bookmark-processor.log
stderr_logfile=/Users/elenafaillace/Library/Logs/supervisor/bookmark-processor.error.log
